# lesson_27/provisioning/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup Timer
Requires=borg-backup.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target